# mining-finance-backend-improved.service
[Unit]
Description=Mining Finance Backend API
Documentation=https://github.com/your-repo/mining-finance
After=network.target docker.service
Wants=docker.service
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/mining-finance/backend
ExecStart=/opt/mining-finance/backend/start-backend.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mining-finance/backend
ReadWritePaths=/var/log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mining-finance-backend

[Install]
WantedBy=multi-user.target



